# conf.d/supervisord.conf

[supervisord]
nodaemon=true

[program:confluence]
command=/opt/atlassian/confluence/bin/start-confluence.sh -fg confluenceuser


